# roles/go-http-server/templates/go-http-server.service.j2
[Unit]
After=local-fs.target network.target

[Service]
Type=simple
ExecStart=/var/www/go-http-server
Restart=on-failure

[Install]
WantedBy=default.target